# faceium.socket
[Unit]
Description=Faceium socket

[Socket]
ListenStream=/run/faceium.sock

[Install]
WantedBy=sockets.target